# nettoyagelogs.timer
[Unit]
Description=Timer pour le nettoyage

[Timer]
OnBootSec=1min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
